# systemd/user/jarvis-tray.service
[Unit]
Description=Jarvis tray indicator
After=graphical-session.target

[Service]
Type=simple
WorkingDirectory=%h/.local/opt/jarvis
Environment=JARVIS_ROOT=%h/.local/opt/jarvis
Environment=JARVIS_STATUS=%h/.local/state/jarvis/status.json
ExecStart=%h/.local/opt/jarvis/.venv/bin/python -m jarvis.tray_app_system
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
