# systemd/user/ollama.service
[Unit]
Description=Ollama local LLM server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/env ollama serve
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
